# openstack-itsl/runner/crontab
# runner crontab — installed for user `runner` by entrypoint.sh
# Staggered slots per CONTRACTS §4: :00/:07/:15/:22 every 30 min per agent.
# Cron is the polling floor; the wake-listener push path is the fast path.
# run-agent.sh holds a per-agent flock, so cron ticks and wakes never overlap.
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin:/app/bin

0,30  * * * * /app/bin/run-agent.sh reb-claude    >> /var/log/runner/cron.log 2>&1
7,37  * * * * /app/bin/run-agent.sh atlas-claude  >> /var/log/runner/cron.log 2>&1
15,45 * * * * /app/bin/run-agent.sh ada-claude    >> /var/log/runner/cron.log 2>&1
22,52 * * * * /app/bin/run-agent.sh marvin-claude >> /var/log/runner/cron.log 2>&1
